# systemd/rpi-cctv.service
[Unit]
Description=Raspberry Pi CCTV Viewer
Documentation=https://github.com/nalbam/rpi
After=graphical.target network-online.target
Wants=network-online.target

[Service]
Type=simple
# Note: Replace 'pi' with your username or adjust WorkingDirectory path
User=pi
WorkingDirectory=/home/pi/rpi/cv2
Environment="DISPLAY=:0"

# Load CCTV credentials from environment file
# Create /etc/rpi/cctv.env with proper permissions (chmod 600)
# Example contents:
#   CCTV_USERNAME=admin
#   CCTV_PASSWORD=your_password
#   CCTV_SERVERS=192.168.1.43:554,192.168.1.44:554
EnvironmentFile=-/etc/rpi/cctv.env

ExecStart=/usr/bin/python3 /home/pi/rpi/cv2/cctv.py
Restart=on-failure
RestartSec=10s

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=rpi-cctv

# Security hardening
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=graphical.target
