# systemd/rpi-kiosk.service
[Unit]
Description=Raspberry Pi Kiosk Mode
Documentation=https://github.com/nalbam/rpi
After=graphical.target network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
Environment="DISPLAY=:0"
ExecStart=/home/pi/start.sh
Restart=on-failure
RestartSec=10s

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=rpi-kiosk

# Security
NoNewPrivileges=true

[Install]
WantedBy=graphical.target
