# systemd/rpi-sonic.service
[Unit]
Description=Raspberry Pi Ultrasonic Distance Sensor
Documentation=https://github.com/nalbam/rpi
After=multi-user.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi
ExecStart=/usr/local/bin/sonic
Restart=on-failure
RestartSec=10s

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=rpi-sonic

# Security hardening
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
